# .local/share/s6/sv/hide-sticky-on-workspace-log/run
#!/bin/execlineb -P
importas -iu ~ HOME
envfile ${~}/.local/share/s6/config/hide-sticky-on-workspace.conf
importas -sCuD "n3 s2000000 T" DIRECTIVES DIRECTIVES
foreground { install -d ${~}/.var/log/hide-sticky-on-workspace }
exec -c s6-log -d3 -b -- ${DIRECTIVES} ${~}/.var/log/hide-sticky-on-workspace
