# .local/share/s6/sv/wireplumber-log/run
#!/bin/execlineb -P
importas -iu ~ HOME
envfile ${~}/.local/share/s6/config/wireplumber.conf
importas -sCuD "n3 s2000000 T" DIRECTIVES DIRECTIVES
foreground { install -d ${~}/.var/log/wireplumber }
exec -c s6-log -d3 -b -- ${DIRECTIVES} ${~}/.var/log/wireplumber
