# .local/share/s6/sv/mako-srv/run
#!/bin/execlineb -P
importas -iu ~ HOME
envfile ${~}/.local/share/s6/config/mako.conf
importas -sCuD "" OPTS OPTS
fdmove -c 2 1
exec mako ${OPTS}
