# .local/share/s6/sv/gnirehtet-srv/run
#!/bin/execlineb -P
importas -iu ~ HOME
envfile ${~}/.local/share/s6/config/gnirehtet.conf
importas -sCuD "" OPTS OPTS
fdmove -c 2 1
exec gnirehtet autorun ${OPTS}
