# .local/share/s6/sv/pause-media-on-headphones-removal-srv/run
#!/bin/execlineb -P
fdmove -c 2 1
foreground { sleep 5 }
exec pause-media-on-headphones-removal
