# .local/share/s6/sv/polkit-gnome-srv/run
#!/bin/execlineb -P
fdmove -c 2 1
exec /usr/lib/polkit-gnome/polkit-gnome-authentication-agent-1
