# .local/share/s6/sv/swayidle-srv/run
#!/bin/execlineb -P
importas -iu ~ HOME
envfile ${~}/.local/share/s6/config/swayidle.conf
importas -uD "" OPTS OPTS
fdmove -c 2 1
exec bash -c "swayidle ${OPTS}"
