# .local/share/s6/sv/wl-clip-persist-log/run
#!/bin/execlineb -P
importas -iu ~ HOME
envfile ${~}/.local/share/s6/config/wl-clip-persist.conf
importas -sCuD "n3 s2000000 T" DIRECTIVES DIRECTIVES
foreground { install -d ${~}/.var/log/wl-clip-persist }
exec -c s6-log -d3 -b -- ${DIRECTIVES} ${~}/.var/log/wl-clip-persist
