# .local/share/s6/sv/wl-clip-persist-srv/run
#!/bin/execlineb -P
fdmove -c 2 1
exec env RUST_LOG=debug wl-clip-persist --clipboard regular --all-mime-type-regex "(?i)^(?!(?:image|audio|video|font|model)/|x-kde-passwordManagerHint).+" --selection-size-limit 524288 --disable-timestamps
