# .local/share/s6/sv/playerctld/up
#!/bin/execlineb -P
exec playerctld daemon
