# .local/share/s6/sv/xdg-desktop-portal-log/run
#!/bin/execlineb -P
importas -iu ~ HOME
envfile ${~}/.local/share/s6/config/xdg-desktop-portal.conf
importas -sCuD "n3 s2000000 T" DIRECTIVES DIRECTIVES
foreground { install -d ${~}/.var/log/xdg-desktop-portal }
exec -c s6-log -d3 -b -- ${DIRECTIVES} ${~}/.var/log/xdg-desktop-portal
